# narwhal/puppet/modules/liferay/pts_motd/files/usr/local/bin/update-motd.d/10-ubuntu-codename
#!/bin/bash

# https://github.com/maxtsepkov/bash_colors/blob/master/bash_colors.sh
uncolorize () { sed -r "s/\x1B\[([0-9]{1,3}((;[0-9]{1,3})*)?)?[m|K]//g"; }

green="\e[1;32m"
nocolor="\e[0m"

sayf() {
  echo -ne "$1";
  echo -e "$nocolor";
}

if [[ $- != *i* ]]
  then
    say() { sayf "$*" ;}

  else
    # do nothing
    say() { echo "$1"; }
fi


ubuntu_codename=$(lsb_release -c -s)

say "Ubuntu codename:$green $ubuntu_codename"
